# supervisord.conf
[program:paid_site]
command = /webapps/paid_site/bin/gunicorn_start     ; Command to start app
user = paid_site                                    ; User to run as
; Where to write log messages
stdout_logfile = /webapps/paid_site/paid_site/logs/gunicorn_supervisor.log  
; Save stderr in the same log
redirect_stderr = true                                          
; Set UTF-8 as default encoding
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8 
